# nodered.service
[Unit]
Description=Node-Red service

[Service]
ExecStart=/usr/local/bin/node-red

[Install]
WantedBy=multi-user.target
