# nodered.service
[Unit]
Description=Node-Red service

[Service]
ExecStart=nodered start

[Install]
WantedBy=multi-user.target
